# life-png.service
[Unit]
Description=Set life.png wallpaper
RefuseManualStart=no
RefuseManualStop=yes
 
[Service]
Type=oneshot
ExecStart=/home/artsi0m/bin/life-png.sh
